# scripts/ironpanel-watchdog.service
[Unit]
Description=IronPanel self-healing watchdog
After=network-online.target ironpanel.service

[Service]
Type=oneshot
ExecStart=/opt/ironpanel/scripts/ironpanel_watchdog.sh
User=root
